# system/launchd/com.threadwork.board-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Autonomous Board sync daemon, made resident (P1 accept gap-closure, task #1641). -->
<!-- Mirrors com.threadwork.board-watcher: KeepAlive + RunAtLoad, run via a wrapper -->
<!-- (run-sync.sh) that cd's to the dashboard root so bun auto-loads dashboard/.env. -->
<!-- SECURITY (#1627): the SUPABASE_SERVICE_ROLE_KEY is NOT in this plist. It is read -->
<!-- by bun from dashboard/.env (gitignored, local-only). No secret value appears here. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.threadwork.board-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/coachstokes/threadwork-dashboard/sync-daemon/launchd/run-sync.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/coachstokes/threadwork-dashboard</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/coachstokes/.threadwork/logs/board-sync/board-sync.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/coachstokes/.threadwork/logs/board-sync/board-sync.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/coachstokes</string>
        <key>PATH</key>
        <string>/Users/coachstokes/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>THREADWORK_VERSION</key>
        <string>v1.0.0</string>
        <key>SYNC_DAEMON_AGENT</key>
        <string>sync-daemon</string>
        <key>SYNC_DAEMON_HEARTBEAT_MS</key>
        <string>30000</string>
        <key>SQLITE_DB_PATH</key>
        <string>/Users/coachstokes/.claude/mcp-servers/task-board/tasks.db</string>
    </dict>
</dict>
</plist>
