# system/launchd/com.threadwork.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.threadwork.watchdog</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/coachstokes/.bun/bin/bun</string>
        <string>run</string>
        <string>/Users/coachstokes/.claude/mcp-servers/task-board/watchdog.ts</string>
    </array>
    <key>StartInterval</key>
    <integer>120</integer>
    <key>StandardOutPath</key>
    <string>/Users/coachstokes/.claude/mcp-servers/task-board/watchdog.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/coachstokes/.claude/mcp-servers/task-board/watchdog.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/coachstokes</string>
        <key>THREADWORK_DEBOUNCE_ENABLED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
