# system/launchd/com.threadwork.bootstrap.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Source: v1/12-j.md "18 write paths" #13 (bootstrap LaunchAgent).
     Rendered via envsubst at install time. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.threadwork.bootstrap</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/coachstokes/.claude/hooks/session-boot.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/coachstokes/.threadwork/logs/bootstrap.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/coachstokes/.threadwork/logs/bootstrap.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/coachstokes</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>THREADWORK_VERSION</key>
        <string>v1.0.0</string>
        <key>THREADWORK_INSTALL_UUID</key>
        <string>REPLACE_WITH_INSTALL_UUID</string>
    </dict>
</dict>
</plist>
